# config/etc/crontabs/root
*/5 * * * * source /var/www/apps/agentex/cronwrapper; python /var/www/apps/agentex/manage.py updateaverages >> /var/www/apps/agentex/update.log
