# config/etc/supervisor.d/processes.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn --workers=4 --bind=127.0.0.1:8080 --access-logfile=- --error-logfile=- core.wsgi:application
directory=/var/www/apps/agentex
user=nginx
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; pid /tmp/nginx.pid;"
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:crond]
command=/usr/sbin/crond -f -L /dev/stdout
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
